# sharelatex/supervisord.ini
[supervisord]
nodaemon=true

[program:chat]
command=/bin/bash -c "/usr/bin/node /var/www/sharelatex/chat/app.js >> /var/log/sharelatex/chat.log 2>&1"

[program:clsi]
command=/bin/bash -c "/usr/bin/node /var/www/sharelatex/clsi/app.js >> /var/log/sharelatex/clsi.log 2>&1"

[program:contacts]
command=/bin/bash -c "/usr/bin/node /var/www/sharelatex/contacts/app.js >> /var/log/sharelatex/contacts 2>&1"

[program:docstore]
command=/bin/bash -c "/usr/bin/node /var/www/sharelatex/docstore/app.js >> /var/log/sharelatex/docstore.log 2>&1"

[program:document-updater]
command=/bin/bash -c "/usr/bin/node /var/www/sharelatex/document-updater/app.js >> /var/log/sharelatex/document-updater.log 2>&1"

[program:filestore]
command=/bin/bash -c "/usr/bin/node /var/www/sharelatex/filestore/app.js /var/log/sharelatex/filestore.log 2>&1"

[program:notifications]
command=/bin/bash -c "/usr/bin/node /var/www/sharelatex/notifications/app.js >> /var/log/sharelatex/notifications.log 2>&1"

[program:real-time]
command=/bin/bash -c "/usr/bin/node /var/www/sharelatex/real-time/app.js >> /var/log/sharelatex/real-time.log 2>&1"

[program:spelling]
command=/bin/bash -c "/usr/bin/node /var/www/sharelatex/spelling/app.js >> /var/log/sharelatex/spelling.log 2>&1"

[program:tags]
command=/bin/bash -c "/usr/bin/node /var/www/sharelatex/tags/app.js >> /var/log/sharelatex/tags.log 2>&1"

[program:track-changes]
command=/bin/bash -c "/usr/bin/node /var/www/sharelatex/track-changes/app.js >> /var/log/sharelatex/track-changes.log 2>&1"

[program:web]
command=/bin/bash -c "/usr/bin/node /var/www/sharelatex/web/app.js >> /var/log/sharelatex/web.log 2>&1"

[program:nginx]
command=/usr/sbin/nginx
stdout_events_enabled=true
stderr_events_enabled=true
